# cfg/crontab
# envmon watchdog
*/1 * * * * /opt/envmon/bin/envmon start
